# files/gerrit.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: gerrit
# REQUIRE: LOGIN
# KEYWORD: shutdown

#
# Add the following line to /etc/rc.conf to enable the Gerrit daemon:
#
# gerrit_enable="YES"
#

. /etc/rc.subr

name="gerrit"
rcvar=`set_rcvar`

load_rc_config "${name}"

: ${gerrit_enable="NO"}
: ${gerrit_home="%%GERRIT_HOME%%"}
: ${gerrit_args=""}
: ${gerrit_java_home="%%JAVA_HOME%%"}
: ${gerrit_user="%%USERS%%"}
: ${gerrit_group="%%GROUPS%%"}
: ${gerrit_site="%%GERRIT_HOME%%/%%GERRIT_SITE%%"}

command="${gerrit_site}/bin/gerrit.sh"
java_cmd="${gerrit_java_home}/bin/java"
procname="GerritCodeReview"
required_files="${java_cmd}"

start_precmd="gerrit_prestart"
start_cmd="gerrit_start"
stop_cmd="gerrit_stop"

PATH="${PATH}:/usr/local/bin"

gerrit_prestart() {
	if [ ! -d "${gerrit_home}" ]; then
		install -d -o "${gerrit_user}" -g "${gerrit_group}" -m 750 "${gerrit_home}"
	fi
	# initialize Gerrit if the site directory is missing
	if [ ! -d "${gerrit_site}" ]; then
		echo "No Gerrit site found. Will initialize Gerrit first..."
		/usr/bin/su - ${gerrit_user} -c "${java_cmd} -jar %%DATADIR%%/gerrit.war init -d ${gerrit_site}"
	fi
}

gerrit_start() {
	${command} start
}

gerrit_stop() {
	${command} stop
}

run_rc_command "$1"
